# extra/iotivitymap.service

[Unit]
Description=
After=network.target

[Service]
ExecStart=$exe
Restart=always

[Install]
WantedBy=default.target

